# python/etc/cron.d/teleinfo
*/1 * * * * root /home/pi/teleinfo/python/teleinfo.sh 
